# Cron/crontab
1 0 * * * root python3 /dailyUpdater.py >> /var/log/cron.log 2>&1

