# Cron/crontab
* * * * * echo "Hello World" >> /var/log/cron.log 2>&1
* * * * * python3 dailyUpdater.py >> /var/log/cron.log 2>&1

